# alpine/supervisord.ini
[supervisord]
nodaemon=true

[program:consumer]
command=python -u consume.py
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
numprocs=%(ENV_PROCESS_NUMBER)s